# quadlet-test/omega-iris.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Iris Voice Assistant
After=omega-infra-pod.service omega-redis.service
SourcePath=%h/.config/containers/systemd/omega-iris.container
RequiresMountsFor=%t/containers
RequiresMountsFor=%h/Documents/Xoe-NovAi/omega-engine/config
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/cache/iris
BindsTo=omega-infra-pod.service
After=omega-infra-pod.service

[X-Container]
Image=localhost/infra_iris:latest
ContainerName=omega-iris
Pod=omega-infra.pod
# Zen 2 Core Steering
# Sovereign Volume
Volume=%h/Documents/Xoe-NovAi/omega-engine/config:/app/config:ro,Z,U
Volume=/media/arcana-novai/omega_library/podman-storage/cache/iris:/cache
Environment=OMEGA_ENTITIES_CONFIG=/app/config/wads/arcana_novai/entities.yaml
Environment=OMEGA_MODELS_CONFIG=/app/config/models.yaml
Environment=REDIS_URL=redis://:omega@localhost:6379
HealthCmd=python3 -c "import urllib.request; urllib.request.urlopen('http://localhost:8080/health')"
Notify=healthy

[Service]
Memory=512M
CPUAffinity=4,6
Restart=always
RestartSec=5
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --userns=keep-id --user 1000 --name omega-iris --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=healthy -d -v %h/Documents/Xoe-NovAi/omega-engine/config:/app/config:ro,Z,U -v /media/arcana-novai/omega_library/podman-storage/cache/iris:/cache --env OMEGA_ENTITIES_CONFIG=/app/config/wads/arcana_novai/entities.yaml --env OMEGA_MODELS_CONFIG=/app/config/models.yaml --env REDIS_URL=redis://:omega@localhost:6379 --health-cmd "python3 -c \"import urllib.request; urllib.request.urlopen('http://localhost:8080/health')" --pod-id-file %t/omega-infra-pod.pod-id localhost/infra_iris:latest

[Install]
WantedBy=default.target
