# quadlet-test/omega-postgres.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Postgres Database
After=omega-infra-pod.service
SourcePath=%h/.config/containers/systemd/omega-postgres.container
RequiresMountsFor=%t/containers
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/volumes/postgres
BindsTo=omega-infra-pod.service
After=omega-infra-pod.service

[X-Container]
Image=docker.io/library/postgres:16-alpine
ContainerName=omega-postgres
Pod=omega-infra.pod
# Zen 2 Core Steering
# Sovereign Volume
Volume=/media/arcana-novai/omega_library/podman-storage/volumes/postgres:/var/lib/postgresql/data
Environment=POSTGRES_USER=omega
Environment=POSTGRES_PASSWORD=omega
Environment=POSTGRES_DB=omega
HealthCmd=pg_isready -U omega
Notify=healthy

[Service]
Memory=512M
CPUAffinity=5,7
Restart=always
RestartSec=5
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --userns=keep-id --user 1000 --name omega-postgres --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=healthy -d -v /media/arcana-novai/omega_library/podman-storage/volumes/postgres:/var/lib/postgresql/data --env POSTGRES_DB=omega --env POSTGRES_PASSWORD=omega --env POSTGRES_USER=omega --health-cmd "pg_isready -U omega" --pod-id-file %t/omega-infra-pod.pod-id docker.io/library/postgres:16-alpine

[Install]
WantedBy=default.target
