# quadlet-test/omega-infra-pod.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Infrastructure Pod
After=network-online.target
SourcePath=/home/arcana-novai/.config/containers/systemd/omega-infra.pod
RequiresMountsFor=%t/containers
Wants=omega-caddy.service
Before=omega-caddy.service
Wants=omega-iris.service
Before=omega-iris.service
Wants=omega-postgres.service
Before=omega-postgres.service
Wants=omega-qdrant.service
Before=omega-qdrant.service
Wants=omega-redis.service
Before=omega-redis.service

[X-Pod]
PodName=omega-infra
# Ports exposed by the pod
PublishPort=127.0.0.1:6333:6333
PublishPort=127.0.0.1:8088:80
PublishPort=127.0.0.1:8080:8080

[Install]
WantedBy=default.target

[Service]
SyslogIdentifier=%N
ExecStart=/usr/bin/podman pod start --pod-id-file=%t/%N.pod-id
ExecStop=/usr/bin/podman pod stop --pod-id-file=%t/%N.pod-id --ignore --time=10
ExecStopPost=/usr/bin/podman pod rm --pod-id-file=%t/%N.pod-id --ignore --force
ExecStartPre=/usr/bin/podman pod create --infra-conmon-pidfile=%t/%N.pid --pod-id-file=%t/%N.pod-id --exit-policy=stop --replace --publish 127.0.0.1:6333:6333 --publish 127.0.0.1:8088:80 --publish 127.0.0.1:8080:8080 --infra-name omega-infra-infra --name omega-infra
Environment=PODMAN_SYSTEMD_UNIT=%n
Type=forking
Restart=on-failure
PIDFile=%t/%N.pid
