# quadlet-test/omega-belial.timer
# 🔱 Omega Belial — Daily Legacy Deep Mining Timer
# AP: AP-OMEGA-BELIAL-TIMER-v1.0.0
#
# Runs the Belial mining service daily at 03:30, staggered from other
# background services to avoid resource contention on the Ryzen 5700U.
#
# Deploy: cp to ~/.config/systemd/user/omega-belial.timer
# Enable: systemctl --user enable --now omega-belial.timer

[Unit]
Description=Omega Belial Daily Legacy Deep Mining Timer
Requires=omega-belial.service

[Timer]
# Daily at 03:30 — staggered from other background tasks
OnCalendar=*-*-* 03:30:00
# Randomized delay to avoid thundering herd if multiple timers fire
RandomizedDelaySec=1800
# Persist missed runs across reboots
Persistent=true

[Install]
WantedBy=default.target