# quadlet-test/omega-caddy.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Caddy Proxy
After=omega-infra-pod.service
SourcePath=/home/arcana-novai/.config/containers/systemd/omega-caddy.container
RequiresMountsFor=%t/containers
RequiresMountsFor=/home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/Caddyfile
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/volumes/caddy/data
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/volumes/caddy/config
BindsTo=omega-infra-pod.service
After=omega-infra-pod.service

[X-Container]
Image=docker.io/library/caddy:2.8-alpine
ContainerName=omega-caddy
Pod=omega-infra.pod
# Zen 2 Core Steering
# Sovereign Volume
Volume=/home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/Caddyfile:/etc/caddy/Caddyfile:ro,Z,U
Volume=/media/arcana-novai/omega_library/podman-storage/volumes/caddy/data:/data
Volume=/media/arcana-novai/omega_library/podman-storage/volumes/caddy/config:/config
HealthCmd=curl -f http://localhost:80/
Notify=healthy

[Service]
Memory=64M
CPUAffinity=9,11
Restart=always
RestartSec=5
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --userns=keep-id --user 1000 --name omega-caddy --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=healthy -d -v /home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/Caddyfile:/etc/caddy/Caddyfile:ro,Z,U -v /media/arcana-novai/omega_library/podman-storage/volumes/caddy/data:/data -v /media/arcana-novai/omega_library/podman-storage/volumes/caddy/config:/config --health-cmd "curl -f http://localhost:80/" --pod-id-file %t/omega-infra-pod.pod-id docker.io/library/caddy:2.8-alpine

[Install]
WantedBy=default.target
