# quadlet-test/omega-redis.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Redis Cache
After=omega-infra-pod.service
SourcePath=/home/arcana-novai/.config/containers/systemd/omega-redis.container
RequiresMountsFor=%t/containers
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/volumes/redis
BindsTo=omega-infra-pod.service
After=omega-infra-pod.service

[X-Container]
Image=docker.io/library/redis:7.4-alpine
ContainerName=omega-redis
Pod=omega-infra.pod
# Zen 2 Core Steering
# Sovereign Volume
Volume=/media/arcana-novai/omega_library/podman-storage/volumes/redis:/data:Z,U
Exec=redis-server --requirepass omega --maxmemory 256mb --maxmemory-policy allkeys-lru --save 300 1 --appendonly yes --appendfsync everysec --no-appendfsync-on-rewrite yes --lazyfree-lazy-eviction yes --lazyfree-lazy-expire yes
HealthCmd=redis-cli -a omega ping
Notify=healthy

[Service]
Memory=256M
CPUAffinity=1,3
Restart=always
RestartSec=5
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --name omega-redis --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=healthy -d -v /media/arcana-novai/omega_library/podman-storage/volumes/redis:/data:Z,U --health-cmd "redis-cli -a omega ping" --pod-id-file %t/omega-infra-pod.pod-id docker.io/library/redis:7.4-alpine redis-server --requirepass omega --maxmemory 256mb --maxmemory-policy allkeys-lru --save 300 1 --appendonly yes --appendfsync everysec --no-appendfsync-on-rewrite yes --lazyfree-lazy-eviction yes --lazyfree-lazy-expire yes

[Install]
WantedBy=default.target
