# quadlet-test/omega-qdrant.service
# Automatically generated by /usr/libexec/podman/quadlet
# 
[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Omega Qdrant Vector Store
After=omega-infra-pod.service
SourcePath=/home/arcana-novai/.config/containers/systemd/omega-qdrant.container
RequiresMountsFor=%t/containers
RequiresMountsFor=/media/arcana-novai/omega_library/podman-storage/volumes/qdrant
RequiresMountsFor=/home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/qdrant_config.yaml
BindsTo=omega-infra-pod.service
After=omega-infra-pod.service

[X-Container]
Image=docker.io/qdrant/qdrant:v1.13.1
ContainerName=omega-qdrant
Pod=omega-infra.pod
# Zen 2 Core Steering
# Sovereign Volume
Volume=/media/arcana-novai/omega_library/podman-storage/volumes/qdrant:/qdrant/storage
Volume=/home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/qdrant_config.yaml:/qdrant/config/config.yaml:ro,Z,U
Environment=QDRANT__SERVICE__GRPC_PORT=6334
Environment=QDRANT__STORAGE__PERFORMANCE__MAX_SEARCH_THREADS=4
Environment=QDRANT__STORAGE__OPTIMIZER__CPU_LIMIT=4
HealthCmd=curl -f http://localhost:6333/healthz
Notify=healthy

[Service]
Memory=1G
CPUAffinity=0,2
Restart=always
RestartSec=5
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --userns=keep-id --user 1000 --name omega-qdrant --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=healthy -d -v /media/arcana-novai/omega_library/podman-storage/volumes/qdrant:/qdrant/storage -v /home/arcana-novai/Documents/Xoe-NovAi/omega-engine/deploy/infra/qdrant_config.yaml:/qdrant/config/config.yaml:ro,Z,U --env QDRANT__SERVICE__GRPC_PORT=6334 --env QDRANT__STORAGE__OPTIMIZER__CPU_LIMIT=4 --env QDRANT__STORAGE__PERFORMANCE__MAX_SEARCH_THREADS=4 --health-cmd "curl -f http://localhost:6333/healthz" --pod-id-file %t/omega-infra-pod.pod-id docker.io/qdrant/qdrant:v1.13.1

[Install]
WantedBy=default.target
